# templates/default/supervisord.filebeat.conf.erb
[program:filebeat]

command = <%=@version%>/filebeat -e -c /var/filebeat.yml

process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup = true
killasgroup = true


stdout_logfile=/tmp/filebeat.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1

stderr_logfile = /tmp/filebeat.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1


[group:filebeat_server]
programs=filebeat




